# deployment/systemd/foundry-agent.service
# Foundry GPU-server agent — reference copy; the installed unit is
# written by `foundry-agent --register` / `--setup-apps` (which also
# trims SupplementaryGroups to groups that exist on the host).
# Pull-only: outbound HTTPS, no listener.
# Docs: docs/DEPLOYMENT.md § GPU Server Runbook · skill: ubuntu-2404-systemd

[Unit]
Description=Foundry GPU server agent
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=foundry-agent
Group=foundry-agent
# Docker socket access is the agent's privilege boundary — docker group,
# not root.
SupplementaryGroups=docker video render
ExecStart=/usr/local/bin/foundry-agent
Environment=FOUNDRY_LOG_FORMAT=json
Restart=on-failure
RestartSec=5
# Sized to the container stop grace period: a deploy/stop task in
# flight finishes or is safely re-dispatched (tasks are idempotent).
TimeoutStopSec=45
# The agent receives only CAP_DAC_OVERRIDE below. SETUID/SETGID/AUDIT_WRITE
# and NET_BIND_SERVICE remain available solely so the setuid-root sudo child
# can initialize and validate/reload nginx through the two sudoers commands.
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_SETUID CAP_SETGID CAP_AUDIT_WRITE CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_DAC_OVERRIDE

# Hardening baseline; relaxations are required by function
# (docs/SECURITY.md § App Publishing):
#   - /etc/foundry-agent: identity config (read-only would break rotation)
#   - /etc/nginx/foundry-apps: agent-managed HTTP/S vhosts
#   - /storage/containers: persistent volume dirs created at deploy
#   - no NoNewPrivileges: the sudoers-scoped `sudo nginx -s reload`
#     needs the setuid transition
#   - ProtectSystem=full, not strict: `nginx -t` writes logs/temp
#     under /var inside this unit's mount namespace
#   - /var/run/docker.sock via SupplementaryGroups
#   - /dev/nvidia*: NVML inventory
#   - CAP_DAC_OVERRIDE: browse project-approved bind mounts even when
#     their container process owns them under an unrelated numeric UID
#   - CAP_SETUID/CAP_SETGID/CAP_AUDIT_WRITE/CAP_NET_BIND_SERVICE are not
#     ambient in the agent; the bounding set retains them for sudo's root
#     transition, audit plugin, and nginx validation of ports 80/443
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/etc/foundry-agent /etc/nginx/foundry-apps /storage/containers

[Install]
WantedBy=multi-user.target
