# deployment/systemd/foundry-backup.service
[Unit]
Description=Foundry MariaDB backup
After=mariadb.service mysql.service

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
ExecStart=/srv/foundry/bin/foundry-backup
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadOnlyPaths=/etc/foundry/mysql-backup.cnf
ReadWritePaths=/srv/foundry/backups/mysql
